# crontab
@reboot python /home/pi/musicbox/musicbox.py &


